# service/install/macos/pkg/scripts/postinstall
#!/bin/bash
#
# ScreenControl Service Postinstall Script
#
# This script runs after the pkg installer copies files.
# It sets up directories, permissions, and starts the service.
#

SERVICE_ID="com.screencontrol.service"
SERVICE_BINARY="/Library/PrivilegedHelperTools/${SERVICE_ID}"
SERVICE_PLIST="/Library/LaunchDaemons/${SERVICE_ID}.plist"
AGENT_ID="com.screencontrol.agent"
AGENT_PLIST="/Library/LaunchAgents/${AGENT_ID}.plist"
CONFIG_DIR="/Library/Application Support/ScreenControl"
LOG_DIR="/Library/Logs/ScreenControl"
APP_PATH="/Applications/ScreenControl.app"

# Log to installer log
log() {
    echo "[ScreenControl Postinstall] $1"
}

log "Starting postinstall..."

# Create directories
log "Creating directories..."
mkdir -p "$CONFIG_DIR"
mkdir -p "$LOG_DIR"

# Set permissions on service binary
log "Setting service permissions..."
chmod 755 "$SERVICE_BINARY"
chown root:wheel "$SERVICE_BINARY"

# Set permissions on plists
chmod 644 "$SERVICE_PLIST"
chown root:wheel "$SERVICE_PLIST"
chmod 644 "$AGENT_PLIST"
chown root:wheel "$AGENT_PLIST"

# Set directory permissions
chmod 755 "$CONFIG_DIR"
chmod 755 "$LOG_DIR"
chown root:wheel "$CONFIG_DIR"
chown root:wheel "$LOG_DIR"

# Create default config if not exists
if [[ ! -f "${CONFIG_DIR}/config.json" ]]; then
    log "Creating default configuration..."
    cat > "${CONFIG_DIR}/config.json" << 'EOF'
{
    "httpPort": 3459,
    "guiBridgePort": 3460,
    "controlServerUrl": "wss://screencontrol.knws.co.uk/ws",
    "agentName": "",
    "autoStart": true,
    "enableLogging": true
}
EOF
    chmod 644 "${CONFIG_DIR}/config.json"
fi

# Load service
log "Loading service..."
launchctl load "$SERVICE_PLIST"

# Wait for service to start
sleep 2

# Verify service is running
if launchctl list | grep -q "$SERVICE_ID" 2>/dev/null; then
    log "Service started successfully"
else
    log "Warning: Service may not have started correctly"
fi

# Load agent for all logged-in GUI users
log "Loading agent for logged-in users..."
for uid in $(dscl . -list /Users UniqueID | awk '$2 >= 500 {print $2}'); do
    # Check if user has a GUI session
    if launchctl print gui/$uid 2>/dev/null | grep -q "state = running"; then
        log "  Loading agent for UID $uid..."
        launchctl bootout gui/$uid/$AGENT_ID 2>/dev/null || true
        launchctl bootstrap gui/$uid "$AGENT_PLIST" 2>/dev/null || true
    fi
done

# Set app permissions (if installed)
if [[ -d "$APP_PATH" ]]; then
    log "Setting application permissions..."
    chmod -R 755 "$APP_PATH"
    chown -R root:wheel "$APP_PATH"

    # Remove quarantine flag
    xattr -rd com.apple.quarantine "$APP_PATH" 2>/dev/null || true
fi

log "Postinstall complete"
log "Service installed at: $SERVICE_BINARY"
log "Agent will auto-start at login for all users"
log "Configuration at: ${CONFIG_DIR}/config.json"
log "Logs at: ${LOG_DIR}/"

exit 0
