# server/bin/security-scan.cron
0 6 * * * root /usr/local/bin/security-scan-autofix.sh >/dev/null 2>&1
